# container-files/etc/supervisor.d/vsphere-import.conf
[program:vsphere-import]
command=vsphere-import
autostart=true
autorestart=true
